# apps/sources-repo/rauldipeas-repo/etc/xdg/autostart/radix-login-sound.desktop
[Desktop Entry]
Type=Application
Name=RaDiX Login Sound
Name[pt_BR]=Som de login do RaDiX
Exec=radix-login-sound
OnlyShowIn=XFCE;