# apps/sources-repo/rauldipeas-repo/etc/rc.local
#!/bin/sh -e
#dpkg-statoverride --remove /var/lib/geoclue
cp /opt/radix-desktop/grub /etc/default/grub
cp /opt/radix-desktop/lsb-release /etc/lsb-release
cp /opt/radix-desktop/os-release /etc/os-release
cp /opt/radix-desktop/menus/* /etc/xdg/menus/
update-grub
exit 0
